# bot.conf
[program:telegrambot]
directory=/home/ubuntu/qorakul_markaz
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stdeer_logfile=/home/ubuntu/qorakul_markaz/logfile_err.log
stdout_logfile=/home/ubuntu/qorakul_markaz/logfile.log
